# azure_vm_setup/systemd/harvey-multi-model-training.service
[Unit]
Description=Harvey Multi-Model Training Data Generator
After=network.target
Documentation=https://github.com/heydividend/harvey

[Service]
Type=oneshot
User=azureuser
Group=azureuser
WorkingDirectory=/opt/harvey
Environment="PATH=/opt/harvey/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="USE_PYMSSQL=true"

# Load environment variables from .env file
EnvironmentFile=/opt/harvey/.env

# Run multi-model training generator with shell wrapper for date substitution
# Generates 800 questions/week (200 per model × 4 models)
ExecStart=/bin/sh -c '/opt/harvey/venv/bin/python /opt/harvey/scripts/multi_model_training_generator.py --category all --questions-per-model 25 --to-database --output /opt/harvey/training_data/multi_model_weekly_$(date +\%Y\%m\%d).json'

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=harvey-multi-model-training

# Resource limits
CPUQuota=80%
MemoryMax=2G
TimeoutStartSec=600

# Restart policy (don't restart for oneshot)
Restart=no

[Install]
WantedBy=multi-user.target
